# opencode-switcher.service
[Unit]
Description=OpenCode Session Switcher
After=graphical-session.target

[Service]
ExecStart=__INSTALL_DIR__/run.sh
Restart=on-failure
RestartSec=3
KillMode=process

[Install]
WantedBy=default.target
