# opencode-switcher.service
[Unit]
Description=OpenCode Session Switcher
After=graphical-session.target

[Service]
ExecStart=/home/hzb/workfiles/opencode-switcher/run.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
